# cake_webdelib/run.cron
#On lance runCrons toutes les minutes (runCrons n'execute que les taches cron en attente)
*/5 * * * * /var/www/webdelib/app/Console/cake -app /var/www/webdelib/app cron > /dev/null 2>&1
